# pianoscript.sh
#!/bin/bash
### BEGIN INIT INFO
# Provides:          piano_learner
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Execute piano_learner python script ...
### END INIT INFO
sudo python3 /home/pi/piano_learner/rpi_main.py